# roles/backend/templates/backend.service.j2
[Unit]
Description = Backend Service

[Service]
User=expense
Environment=DB_HOST="{{DB_HOST}}"
ExecStart=/bin/node /app/index.js
SyslogIdentifier=backend

[Install]
WantedBy=multi-user.target